# root/etc/periodic/weekly/dehydrated
#!/bin/bash
chown -R ${UID}:${GID} /data
# Run dehydrated
su-exec ${UID}:${GID} /opt/dehydrated/dehydrated --config /data/config --cron --keep-going
