# root/etc/periodic/weekly/dehydrated
#!/bin/bash
chown -R ${UID}:${GID} /etc/dehydrated /certs /var/www/dehydrated
su-exec ${UID}:${GID} /dehydrated/dehydrated -c
